# docker/cron/sqlite.crontab
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# Every 30 minutes: clean up stale reset tokens and unverified users
*/30 * * * * /bin/sh /bookshop-server/docker/cron/jobs/cleanup.sh

# Daily at 02:00: create SQLite backup
0 2 * * * /bin/sh /bookshop-server/docker/cron/jobs/backup_sqlite.sh
